# debian/extra/units/systemd-vconsole-setup.service.d/setupcon-noop.conf
[Unit]
# Prefer using console-setup if that is installed
ConditionPathExists=!/bin/setupcon
